# config/bash/profile.d/secrets.sh.example
#!/usr/bin/env bash
export SOME_API_TOKEN=aabbccddeeff
